# my-project/cloud-drive.service
[Unit]
Description=Cloud Drive API
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser
Environment="CLOUD_DRIVE_ROOT=/home/azureuser/cloud-drive"
ExecStart=/usr/bin/python3 /home/azureuser/bin/cloud_drive_server.py --host 127.0.0.1 --port 8787 --root /home/azureuser/cloud-drive
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloud-drive

[Install]
WantedBy=multi-user.target