# my-project/cloudflared.service
[Unit]
Description=Cloudflare Named Tunnel for Gateway
After=network.target network-online.target cloud-gateway.service
Wants=network-online.target cloud-gateway.service

[Service]
Type=simple
User=azureuser
ExecStart=/usr/bin/cloudflared --config /home/azureuser/.cloudflared/config.yml tunnel run
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared

[Install]
WantedBy=multi-user.target
