# my-project/cloudflared.quick-tunnel.service
[Unit]
Description=Cloudflare Quick Tunnel for Gateway
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=azureuser
ExecStart=/usr/bin/cloudflared tunnel --url http://localhost:8080
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
