# my-project/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for Gateway (Token Managed)
After=network.target network-online.target cloud-gateway.service
Wants=network-online.target cloud-gateway.service

[Service]
Type=simple
User=azureuser
EnvironmentFile=-/home/azureuser/.cloudflared/cloudflared-token.env
ExecStart=/usr/bin/cloudflared tunnel run --token ${TUNNEL_TOKEN}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared

[Install]
WantedBy=multi-user.target
