# my-project/cpolar.service
[Unit]
Description=Cpolar HTTP Tunnel for LobeChat
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/azureuser/bin/cpolar http 3210 -authtoken cdec8507-3ce7-40a1-9be2-d2406c3cfaae -region cn
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cpolar

[Install]
WantedBy=multi-user.target
